# deploy/sms-gateway.service
[Unit]
Description=SMS Gateway API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sms-gateway
EnvironmentFile=/opt/sms-gateway/.env
ExecStart=/opt/sms-gateway/.venv/bin/python -m api
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
