# deploy/sms-gateway-admin.service
[Unit]
Description=SMS Gateway Admin
After=network-online.target sms-gateway.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sms-gateway
EnvironmentFile=/opt/sms-gateway/.env
Environment=APP_ROLE=admin
Environment=APP_PORT=8089
ExecStart=/opt/sms-gateway/.venv/bin/python -m api
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
